# wargame/src/server/rootfs/etc/service/nginx/run
#!/bin/sh -e

# pipe stderr to stdout and run nginx omiting ENV vars to avoid security leaks
exec 2>&1
exec env - PATH=$PATH nginx -g 'daemon off;'